# yomypopcorn-scanner.service.j2
[Unit]
Description=yomypopcorn-scanner

[Service]
Restart=always
ExecStart={{ INSTALL_PATH }}/scanner/bin/cli.js --debug
